# templates/default/rhel7-eja-dicom.service.erb
[Unit]
Description=Start dicom listener 
After=network.target

[Service]
Type=forking
ExecStart=<%= @startup_script %>
TimeoutStartSec=0

[Install]
WantedBy=default.target

